# files/systemd/system/nightly-shutdown.service
[Unit]
Description=Nightly Shutdown Service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff